# debian/beacon-gpx.socket
[Unit]
Description=Beacon GPX converter
Before=nginx.service

[Socket]
ListenStream=%t/beacon/gpx.socket

[Install]
WantedBy=nginx.service
