# debian/beacon-api.socket
[Unit]
Description=Beacon API
Before=nginx.service

[Socket]
ListenStream=%t/beacon/api.socket

[Install]
WantedBy=nginx.service
